# launchd/plists/com.snackdriven.auto-organize-screengrabs.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.snackdriven.auto-organize-screengrabs</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/kayla/.local/bin/auto-organize-screengrabs.sh</string>
    </array>
    <key>StartInterval</key>
    <integer>3600</integer>
    <key>RunAtLoad</key>
    <false/>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/kayla</string>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/kayla/.local/bin/.auto-organize-screengrabs.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/kayla/.local/bin/.auto-organize-screengrabs.err</string>
</dict>
</plist>
